# docker/crontab/entrypoint.sh
#!/usr/bin/bash

/sbin/tini -- /docker-entrypoint crond -f