# docker/crontab/entrypoint.sh
#!/usr/bin/env bash

/sbin/tini -- /docker-entrypoint crond -f